# supervisor-anki-server.conf
[program:anki-server]
command=/usr/local/bin/paster serve production.ini
directory=/var/anki-server
user=www-data
autostart=true
autorestart=true
redirect_stderr=true

