# supervisor-anki-server.conf
[program:anki-server]
command=/usr/local/bin/paster serve %(here)s/production.ini
autostart=true
autorestart=true
user=www-data
redirect_stderr=true

